# loki/rootfs/etc/services.d/loki/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Loki
# Runs Loki
# ==============================================================================

bashio::log.info 'Starting Loki...'
loki_config='/etc/loki/default-config.yaml'

if bashio::config.exists 'config_path'; then
    loki_config=$(bashio::config 'config_path')
    bashio::log.info "Using config at ${loki_config}"
else
    bashio::log.info "Using default config"
fi

days_to_keep=30
if bashio::config.exists 'days_to_keep'; then
    days_to_keep=$(bashio::config 'days_to_keep')
fi
# Subtract one because Loki keeps one more index period then we say
# https://grafana.com/docs/loki/latest/operations/storage/table-manager/#retention
days_to_keep=$((days_to_keep - 1))
retention_period="$((days_to_keep * 24))h"
bashio::log.info "Retention period set to ${retention_period}"
export "RETENTION_PERIOD=${retention_period}"

case "$(bashio::config 'log_level')" in \
    error)	log_level='error' ;; \
    warning)	log_level='warn' ;; \
    debug)	log_level='debug' ;; \
    *)		log_level='info' ;; \
esac;
bashio::log.info "Loki log level set to ${log_level}"

loki_args=(
    "-config.expand-env=true"
    "-config.file=${loki_config}"
    "-server.http-listen-address=127.0.0.1"
    "-server.http-listen-port=80" 
    "-log.level=${log_level}"
)
if [ "${log_level}" == "debug" ]; then
    bashio::log.debug "Logging full config on startup for debugging..."
    loki_args+=("-print-config-stderr=true")
fi

bashio::log.info "Handing over control to Loki..."
loki "${loki_args[@]}"
